# support/mqtt/wardrive-mqtt.service
[Unit]
Description=MeshCore Coverage Bot (MQTT)

[Service]
WorkingDirectory=/home/ubuntu/mesh-map/support/mqtt
ExecStart=/home/ubuntu/venv/bin/python3 -u /home/ubuntu/mesh-map/support/mqtt/wardrive-mqtt.py
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/ubuntu/venv/bin:/usr/local/bin:/usr/bin"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
User=root
Group=root

[Install]
WantedBy=multi-user.target